# com.mill-cache.gc.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mill-cache.gc</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/tunguyen/bin/mill-cache</string>
        <string>gc</string>
    </array>

    <key>StartInterval</key>
    <integer>3600</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/mill-cache-gc.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/mill-cache-gc.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>DEBUG</key>
        <string>0</string>
    </dict>
</dict>
</plist>
